# init/systemd/egress-proxy.service
[Unit]
Description=Agents Secure Binding Egress Proxy Service
After=network.target

[Service]
EnvironmentFile=/etc/agents-secure-binding/environment
Environment="ASB_LOG_LEVEL=debug"
ExecStart=/usr/bin/agents-secure-binding-egress-proxy --port=3128
Restart=always
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target
